# honeypot.service
[Unit]
Description=OpenAI API Honeypot
Documentation=https://github.com/your-repo/honeypot
After=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/root/Honey-AI
Environment="PATH=/root/Honey-AI/venv/bin:/usr/bin:/bin"
EnvironmentFile=-/root/Honey-AI/.env
ExecStart=/root/Honey-AI/venv/bin/python /root/Honey-AI/main.py
ExecReload=/bin/kill -HUP $MAINPID
Restart=always
RestartSec=5

# Security hardening
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=no
ReadWritePaths=/root/Honey-AI

# Resource limits
LimitNOFILE=65535
LimitNPROC=4096

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=honeypot

[Install]
WantedBy=multi-user.target
